# config/cron/crontab
# SUPERCRONIC CRONTAB - V3 Consolidated (2026-03-09)
# Format: minute hour day month weekday command
# Container: docker-compose.worker.yml (Supercronic + BullMQ)

# =====================================================
# SEO LINK TRACKING (Maillage interne)
# =====================================================

# SEO - Aggregation metriques quotidiennes (2h30)
30 2 * * * curl -s -X POST http://localhost:3000/api/seo/aggregate >> /logs/seo-aggregate.log 2>&1

# SEO - Nettoyage donnees > 90 jours (lundis 3h30)
30 3 * * 1 curl -s -X POST "http://localhost:3000/api/seo/cleanup?daysToKeep=90" >> /logs/seo-cleanup.log 2>&1

# =====================================================
# MAINTENANCE SYSTEME
# =====================================================

# Cache - Nettoyage Redis (5h)
0 5 * * * node /app/backend/dist/scripts/cleanup-redis-cache.js >> /logs/redis-cleanup.log 2>&1

# Payment tunnel - 10 tests read-only (toutes les 15 min)
*/15 * * * * /app/scripts/test-payment-tunnel.sh >> /logs/payment-tunnel.log 2>&1 || echo "[$(date)] Payment tunnel FAILED" >> /logs/health.log

# BullMQ - Stats quotidiennes (6h)
0 6 * * * node /app/backend/dist/scripts/bullmq-stats.js >> /logs/bullmq-stats.log 2>&1

# Queue Health - Toutes les queues BullMQ (15 min)
*/15 * * * * /app/scripts/cron/queue-health.sh >> /logs/queue-health.log 2>&1

# Sitemap Freshness - Verification hebdo (lundi 4h)
0 4 * * 1 /app/scripts/cron/sitemap-freshness.sh >> /logs/sitemap-freshness.log 2>&1

# =====================================================
# BACKUP & LOGS
# =====================================================

# Backup - Dump Redis (1h)
0 1 * * * redis-cli --rdb /backups/redis-$(date +\%Y\%m\%d).rdb 2>&1

# Logs - Rotation (dimanches 23h)
0 23 * * 0 find /logs -name "*.log" -mtime +7 -delete
