# scripts/reindex_two_pass.cron
# Install with: crontab /opt/automecanik/app/scripts/reindex_two_pass.cron
# Runs every night at 02:15 server time.
15 2 * * * /opt/automecanik/app/scripts/reindex_two_pass_cron.sh >> /opt/automecanik/rag/logs/reindex-cron/cron_runner.log 2>&1
